# supervisord.ini
[supervisord]
nodaemon=true

[program:monetdbd]
command=monetdbd start /var/monetdb5/dbfarm
autorestart=false
user=monetdb
priority=1
exitcodes=0

[program:monetdb]
command=monetdb start db
autorestart=false
user=monetdb
priority=2
exitcodes=1
