# supervisord.conf
[supervisord]
nodaemon=true

[program:mserver5]
command=/usr/bin/mserver5 --dbpath=/var/dmonetdb5/db --daemon=yes --set embedded_r=true >> /var/dmonetdb5/db/monetdb.log
autorestart=true
