# supervisord.conf
[supervisord]
nodaemon=true

[program:mserver5]
command=sh /root/start-monetdb.sh db
autorestart=true
