# supervisord.conf
[supervisord]
nodaemon=true

[program:monetdbd]
command=/usr/bin/monetdbd start /var/monetdb5/dbfarm
autorestart=true

[program:monetdb]
command=/usr/bin/monetdb start db
autorestart=true
