# supervisord.ini
[supervisord]
nodaemon=true

[program:monetdbd]
command=monetdbd start /var/monetdb5/dbfarm
autorestart=true
user=monetdb

[program:monetdb]
command=monetdb start db
autorestart=true
user=monetdb
