# configs/supervisord.conf
[supervisord]
nodaemon=true

[program:monetdbd]
command=monetdbd start /var/monetdb5/dbfarm
autorestart=false
user=monetdb
priority=1
exitcodes=0
startsecs=0
